# docker/supervisor/conf.d/laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/html
command=/usr/local/bin/php /var/www/html/artisan queue:work --queue=default --sleep=3 --tries=3
user=root
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
loglevel=blather
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=14
stdout_logfile=/var/log/supervisor/laravel-worker.log
stopwaitsecs=86400
stopasgroup=true
killasgroup=true

; [program:laravel-schedule-run]
; process_name=%(program_name)s_%(process_num)02d
; command=/bin/bash -c "while [ true ]; do (php /var/www/html/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
; autostart=true
; autorestart=true
; user=root
; numprocs=1
; redirect_stderr=true
; stdout_logfile=/var/log/supervisor/laravel-schedule.log
; stopwaitsecs=60
